# systemd-service/extboot-sync.service
[Unit]
Description=Extboot sync (stub)

[Service]
Type=oneshot
ExecStart=/usr/bin/touch /tmp/extboot-sync
